# supervisor.d/go-client.ini

[program:gq-client]
command=/usr/local/bin/gq-client -s 47.244.243.77 -l 1984 -c /Users/daipei/etc/gqclient.json
autostart=true
autorestart=true
stdout_logfile=/Users/daipei/.supervisor/logs/gqclient.info
stderr_logfile=/Users/daipei/.supervisor/logs/gqclient.err
