# supervisor.d/sync.ini

[program:dotfiles-sync]
command=/Users/daipei/dotfiles/sync.sh
autostart=true
autorestart=true
stdout_logfile=/Users/daipei/.supervisor/dotfiles-sync.info
stderr_logfile=/Users/daipei/.supervisor/dotfiles-sync.err
